# docker/templates/crontab
*/1 * * * *     root    /var/www/vegadns2/bin/update-data.sh
